# scripts/mac/com.constitution-memorizer.serve.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent template for Constitution Memorizer UI.
  Install via: scripts/mac/install-serve-agent.sh
  Placeholders __REPO_ROOT__ and __VENV_PYTHON__ are replaced by the installer.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.constitution-memorizer.serve</string>
  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>
  <key>ProgramArguments</key>
  <array>
    <string>__VENV_PYTHON__</string>
    <string>-m</string>
    <string>constitution_memorizer.cli</string>
    <string>serve</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>8001</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/constitution-memorizer-serve.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/constitution-memorizer-serve.err.log</string>
</dict>
</plist>
