# scripts/mac/com.constitution-memorizer.reminders.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Daily study reminder LaunchAgent (ntfy).
  Install via: scripts/mac/install-reminders-agent.sh
  Placeholders: __REPO_ROOT__, __VENV_PYTHON__, __HOME__, __NTFY_TOPIC__, __NTFY_SERVER__, __NTFY_TOKEN__
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.constitution-memorizer.reminders</string>
  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>
  <key>ProgramArguments</key>
  <array>
    <string>__VENV_PYTHON__</string>
    <string>-m</string>
    <string>constitution_memorizer.cli</string>
    <string>send-reminders</string>
    <string>--channel</string>
    <string>ntfy</string>
    <string>--host-url</string>
    <string>http://127.0.0.1:8001</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>NTFY_TOPIC</key>
    <string>__NTFY_TOPIC__</string>
    <key>NTFY_SERVER</key>
    <string>__NTFY_SERVER__</string>
    <key>NTFY_TOKEN</key>
    <string>__NTFY_TOKEN__</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/constitution-memorizer-reminders.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/constitution-memorizer-reminders.err.log</string>
</dict>
</plist>
